# ops/systemd/metaos-cleanup.timer
[Unit]
Description=Run METAOS cleanup

[Timer]
OnBootSec=10m
OnUnitActiveSec=30m
Persistent=true

[Install]
WantedBy=timers.target
