# ops/systemd/metaos-cleanup.service
[Unit]
Description=METAOS cleanup

[Service]
Type=oneshot
ExecStart=/opt/metaos/app/ops/cleanup-artifacts.sh
